# roles/setup-pxe/templates/tftp.service.j2
[Unit]
Description=Tftp Server
Requires=tftp.socket
Documentation=man:in.tftpd

[Service]
ExecStart=/usr/sbin/in.tftpd -c -p -s {{ TFTP_DIR }}
StandardInput=socket

[Install]
Also=tftp.socket